# backup_rsync/rsync.target
[Unit]
Description=Minute Timer Targer
StopWhenUnneeded=yes